# deploy/systemd/neo-node.service
[Unit]
Description=Neo C++ Blockchain Node
Documentation=https://github.com/neo-project/neo-cpp
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=neo
Group=neo
WorkingDirectory=/opt/neo-cpp

# Start command with configuration
ExecStart=/opt/neo-cpp/bin/neo-node \
    --config /etc/neo/config.json \
    --datadir /var/lib/neo \
    --rpc \
    --daemon

# Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=300
KillMode=mixed
KillSignal=SIGTERM

# Restart configuration
Restart=always
RestartSec=30
StartLimitInterval=600
StartLimitBurst=5

# Process management
LimitNOFILE=65535
LimitNPROC=4096
MemoryLimit=8G
CPUQuota=400%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Filesystem access
ReadWritePaths=/var/lib/neo /var/log/neo
ReadOnlyPaths=/etc/neo

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neo-node

# Environment
Environment="NEO_LOG_LEVEL=info"
Environment="NEO_NETWORK=mainnet"
EnvironmentFile=-/etc/neo/environment

[Install]
WantedBy=multi-user.target